# scripts/home_server/systemd/cdx-trade-nightly.service
[Unit]
Description=cdx-trade nightly pipeline
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/<user>/cdx-trade
ExecStart=/home/<user>/cdx-trade/scripts/home_server/run_nightly.sh
Environment=CONFIG_PATH=configs/etf_core_1k.json
Environment=KEEP_DAYS=30
Environment=KEEP_MAX_RUNS=120
Environment=MAX_PENDING_KEEP=1
Environment=MIN_PENDING_SCORE=0.00
Environment=PENDING_SCORE_LOOKBACK_DAYS=180
EnvironmentFile=-/home/<user>/cdx-trade/.env

[Install]
WantedBy=multi-user.target
