# scripts/home_server/systemd/cdx-trade-nightly.timer
[Unit]
Description=Run cdx-trade nightly pipeline daily

[Timer]
OnCalendar=Mon..Fri 21:15
Persistent=true

[Install]
WantedBy=timers.target
